# image/templates/files/install.xml
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Copyright 2024 Oxide Computer Company
-->

<!--
    We must use the service name "system/install" here because SMF
    requires either console-login to be able to come up, or the instance
    "svc:/system/install:default".  If one of those cannot come up,
    the system goes into maintenance mode.
-->
<service_bundle type='manifest' name='site-install'>
<service name='system/install' type='service' version='1'>
    <create_default_instance enabled='true' />

    <!-- Must be able to access /tmp. -->
    <dependency
        name='single-user'
        grouping='require_all'
        restart_on='none'
        type='service'>
        <service_fmri value='svc:/milestone/single-user' />
    </dependency>

    <exec_method
        type='method'
        name='start'
        exec='/usr/lib/captive &gt;/dev/console &lt;/dev/console 2&gt;&amp;1'
        timeout_seconds='0' />

    <exec_method
        type='method'
        name='stop'
        exec=':kill -9'
        timeout_seconds='3' />

    <property_group name='startd' type='framework'>
        <propval name='duration' type='astring' value='child' />
        <propval name='ignore_error' type='astring' value='core,signal' />

        <!-- this makes our process a session leader on /dev/console: -->
        <propval name='need_session' type='boolean' value='true' />
    </property_group>

    <stability value='Unstable' />

    <template>
        <common_name>
            <loctext xml:lang='C'>
                captive console (replaces console login)
            </loctext>
        </common_name>
    </template>
</service>
</service_bundle>
