# image/templates/files/netstack.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='profile' name='helios'
  xmlns:xi='http://www.w3.org/2003/XInclude' >
  <service name='network/ntp' version='1' type='service'>
    <instance name='default' enabled='true' />
  </service>
  <service name='system/illumos/mg-ddm' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
  <service name='system/illumos/opte' version='1' type='service'>
    <instance name='default' enabled='false' />
  </service>
</service_bundle>
