# image/templates/files/site-install.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    Copyright 2024 Oxide Computer Company
-->

<service_bundle type='profile' name='helios'>
    <service name='system/hotplug' version='1' type='service'>
        <instance name='default' enabled='true' />
    </service>

    <!-- See: https://www.illumos.org/issues/14006 -->
    <service name='network/routing/route' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='network/routing-setup' version='1' type='service'>
        <instance name='default'>
            <property_group name='routeadm' type='application'>
                <propval name='ipv4-routing-set' type='boolean' value='true' />
            </property_group>
        </instance>
    </service>

    <!-- Disable console login service so that boot menu can take over: -->
    <service name='system/console-login' version='1' type='service'>
        <instance name='default' enabled='false'/>
    </service>
</service_bundle>
